# openclaw-pi-dashboard.service
[Unit]
Description=OpenClaw Pi Dashboard
After=network.target

[Service]
Type=simple
User=imageidentity
WorkingDirectory=/home/imageidentity/.openclaw/workspace/projects/OpenClaw-Pi-Dashboard
Environment="PATH=/home/imageidentity/.openclaw/workspace/projects/OpenClaw-Pi-Dashboard/venv/bin"
Environment="DASHBOARD_CONFIG=/home/imageidentity/.openclaw/workspace/projects/OpenClaw-Pi-Dashboard/config.json"
ExecStart=/home/imageidentity/.openclaw/workspace/projects/OpenClaw-Pi-Dashboard/venv/bin/python /home/imageidentity/.openclaw/workspace/projects/OpenClaw-Pi-Dashboard/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
