# openclaw-pi-dashboard.service
[Unit]
Description=OpenClaw Pi Dashboard
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/OpenClaw-Pi-Dashboard
Environment="PATH=/home/pi/OpenClaw-Pi-Dashboard/venv/bin"
Environment="DASHBOARD_CONFIG=/home/pi/OpenClaw-Pi-Dashboard/config.json"
ExecStart=/home/pi/OpenClaw-Pi-Dashboard/venv/bin/python /home/pi/OpenClaw-Pi-Dashboard/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
